# deploy/kalshi-daemon.service
# Kalshi Unified Trading Daemon (Phase 1)
#
# Persistent process. Replaces:
#   - kalshi-bot.service + kalshi-bot.timer (the 2-minute oneshot)
#   - kalshi-weather-daemon.service (weather MM daemon)
# Both of those should be STOPPED and DISABLED before this starts.
#
# Architecture:
#   - Main thread: Scheduler running cycle (60s), kv_cleanup (1h), health_log (5m)
#   - Poller threads: METARPoller (30s). Phase 2 adds NWS/NBM/HRRR/MADIS/AFD.
#   - Shared SQLite connection in WAL mode; writes serialized by DB_WRITE_LOCK.
#
# Install:
#   cp deploy/kalshi-daemon.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl disable --now kalshi-bot.timer
#   systemctl disable --now kalshi-weather-daemon  # if present
#   systemctl enable --now kalshi-daemon
#
# Logs: tail -f /home/kalshi/autoagent/daemon.log

[Unit]
Description=Kalshi Unified Trading Daemon
After=network-online.target
Wants=network-online.target

# Belt-and-suspenders — fail to start if the old units are still active.
# They should be disabled before enabling this unit.
Conflicts=kalshi-bot.service kalshi-bot.timer kalshi-weather-daemon.service

# 2026-05-04: rate-limiter directives moved here from [Service]. Modern
# systemd (>= 230) requires StartLimit* in [Unit], not [Service]. Prior
# placement caused "Unknown key name 'StartLimitIntervalSec'" warnings
# and the directives were silently ignored — defaults of 5 starts /
# 10s applied instead. Generous limits below give the daemon room to
# crash-loop a few times before systemd gives up entirely.
StartLimitIntervalSec=600
StartLimitBurst=10

[Service]
Type=simple
User=kalshi
Group=kalshi
WorkingDirectory=/home/kalshi/autoagent
ExecStart=/usr/bin/python3 -m bot.daemon.main

# 2026-05-04: changed Restart=on-failure → Restart=always after the
# 9h outage on May 3. SIGKILL-after-TimeoutStopSec IS technically
# "result=timeout" (a failure), but systemd's cgroup teardown
# encountered "Failed to kill control group ... Invalid argument"
# which left the unit in a state where on-failure didn't auto-recover.
# Restart=always restarts on EVERY exit — safer default for an
# always-on trading daemon.
Restart=always
RestartSec=30

# 2026-05-04: lowered TimeoutStopSec 360 → 90.
# Original 360 was sized for mos_materializer's ~300s worst case, but
# during a routine deploy hitting that limit means SIGKILL + a
# multi-hour auto-restart failure window (see May 3 outage). 90s lets
# a deploy fail-fast if the daemon is wedged so systemd can recycle it
# cleanly. The scheduler's on_stop hooks should complete well under
# this; any task taking >90s on shutdown is a separate bug to fix in
# code, not paper over with longer timeouts.
TimeoutStopSec=90

# Environment
EnvironmentFile=/home/kalshi/autoagent/.env
Environment=PYTHONUNBUFFERED=1
Environment=KALSHI_LOG_LEVEL=INFO

# Logging to file (append, matches the oneshot pattern so dashboards
# don't have to learn two file names).
StandardOutput=append:/home/kalshi/autoagent/daemon.log
StandardError=append:/home/kalshi/autoagent/daemon.log

# Resource limits. 30 days × many threads needs more than the oneshot's
# nothing-in-particular limit, but well under VPS capacity.
MemoryMax=512M
CPUQuota=50%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/kalshi/autoagent
PrivateTmp=true

[Install]
WantedBy=multi-user.target
